# cache1/supervisor/trafficserver.conf
[program:traffic_cop]
command=/opt/trafficserver/bin/traffic_cop
numprocs=1
autostart=true
autorestart=true
stopsignal=QUIT
startsecs=1
