# cache/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
numprocs=1
autostart=true
autorestart=true
stopsignal=QUIT
startsecs=1
