# cache/supervisor/trafficserver-first.conf
[program:trafficserver_first]
command=/opt/trafficserver-first/bin/traffic_cop
numprocs=1
autostart=true
autorestart=true
stopsignal=QUIT
startsecs=1
environment=TS_ROOT="/opt/trafficserver-first"
