# cache/supervisor/trafficserver-second.conf
[program:trafficserver_second]
command=/opt/trafficserver-second/bin/traffic_cop
numprocs=1
autostart=true
autorestart=true
stopsignal=QUIT
startsecs=1
environment=TS_ROOT="/opt/trafficserver-second"
